# deploy/systemd/blog-scan-orphan-files.service
[Unit]
Description=Blog scan orphan local files
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/blog
ExecStart=/usr/bin/docker compose -f deploy/docker-compose.yml -f deploy/docker-compose.prod.yml exec -T backend uv run python -m app.cli cleanup-orphan-files --limit 1000
